# refacto/clean-mongo-log.service
[Unit]
Description=Clean MongoDB Logs

[Service]
Type=oneshot
ExecStart=/var/www/refacto/clean-mongo-log.sh
