# refacto/refacto-updater.service
[Unit]
Description=Refacto Updater

[Service]
Type=oneshot
ExecStart=/var/www/refacto/update.sh
After=network-online.target
