# refacto/refacto-updater.service
[Unit]
Description=Refacto Updater
After=network.target

[Service]
Type=oneshot
ExecStart=/var/www/refacto/update.sh
