# sequence/sequence-updater.service
[Unit]
Description=Sequence Diagram Updater

[Service]
Type=oneshot
ExecStart=/var/www/sequence/update.sh
After=network-online.target
